# home/.config/systemd/user/autoforge.service
[Unit]
Description=Autoforge Service
After=network.target

[Service]
Type=simple
# autoforge installed via: bun add -g autoforge-ai
ExecStart=%h/.bun/bin/autoforge
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
