# home/.config/systemd/user/mic-ttgo.service
[Unit]
Description=Mic Mute/Unmute → TTGO Display
After=default.target

[Service]
# uv run reads the PEP 723 inline metadata in ttgo2.py and auto-manages
# an isolated venv with pyserial. No pyenv or manual pip install needed.
# uv is managed by mise; shim path is stable after 'mise install'.
ExecStart=%h/.local/share/mise/shims/uv run %h/apps/ttgo2.py
WorkingDirectory=%h/apps
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
