# home/.config/systemd/user/mic-ttgo.service
[Unit]
Description=Mic Mute/Unmute → TTGO Display
After=default.target

[Service]
# Python path via pyenv — update version if you change pyenv default
ExecStart=%h/.pyenv/versions/3.13.7/bin/python3 %h/apps/ttgo2.py
WorkingDirectory=%h/apps
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5

# Uncomment if PulseAudio/Pipewire needs to be running first:
# Environment=PULSE_SERVER=unix:/run/user/%U/pulse/native

[Install]
WantedBy=default.target
